# com.jeff.fable-self-improve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- DRAFT — load only after Jeff approves: launchctl load ~/grokgo/com.jeff.fable-self-improve.plist
     Runs the reflection pass daily at 09:00. The script self-throttles (no-op under MIN_NEW new
     transcripts), so a daily schedule is cheap. Unload to disable. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jeff.fable-self-improve</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/rentamac/grokgo/self-improve.sh</string>
    <string>--once</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/rentamac/grokgo/logs/self-improve.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/rentamac/grokgo/logs/self-improve.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
